# config_files/home/tv/.config/systemd/user/casting.service
[Unit]
Description=Run casting server

[Service]
ExecStart=/home/tv/casting_server.py
Restart=always

[Install]
WantedBy=default.target